# .config/systemd/user/batt-check.timer
[Unit]
Description=low battery alert service

[Timer]
OnCalendar=*:00/5

[Install]
WantedBy=default.target
