# .config/systemd/user/batt-check.service
[Unit]
Description=low battery alert service

[Service]
ExecStart=$HOME/.local/bin/batt-check

[Install]
WantedBy=default.target
