# etc/cron.hourly/logmail
#!/bin/sh
journalctl --since=-1h --priority=err --no-pager --quiet | grep -v '^-- No entries --$'
systemctl --state=failed --no-pager --no-legend
