# etc/profile.d/stdman.sh
export MANPATH="/usr/local/share/man:$MANPATH"
