# etc/systemd/system/ddns.timer
[Unit]
Description=update Linode DNS IP address

[Timer]
OnCalendar=*:00/5

[Install]
WantedBy=timers.target
