# etc/systemd/system/ddns.timer
[Unit]
Description=update Linode DNS IP address

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
