# etc/systemd/system/ddns.service
[Unit]
Description=update Linode DNS IP address

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ddns
