# etc/systemd/system/getty@.service.d/nohints.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --nohints --noclear %I $TERM
